# ansible/roles/mailhog/templates/mailhog.service.j2
[Unit]
Description=MailHog Email Catcher
After=network.target

[Service]
User={{ ansible_user_id }}
Group={{ ansible_user_id }}
EnvironmentFile=/etc/mailhog.env
ExecStart=/usr/local/bin/mailhog
Restart=always

[Install]
WantedBy=multi-user.target
